# volume-control-panel@.service
[Unit]
Description=Controls volume from panel at %I

[Service]
Type=oneshot
ExecStart=/usr/local/bin/log-activation %I
