# volume-control-panel@.service
[Unit]
Description=Controls volume from panel at %I

[Service]
Type=simple
ExecStart=/usr/bin/volume-control-panel-interpreter %I
